# reflink_dedupe.cron
0 2 1 * * root /usr/local/bin/reflink_dedupe --schedule 2>&1
